# config/systemd/user/gnk-battery-monitor.service
[Unit]
Description=GNK Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/gnk/bin/gnk-battery-monitor
Environment=DISPLAY=:0
